# go/packaging/root/usr/share/stalker/init/stalkerd.conf
start on (filesystem and net-device-up and remote-filesystem)

chdir /var/lib/stalker
console log

kill signal SIGINT
kill timeout 20

script 
    /usr/bin/stalkerd
end script
